# claude/.claude/plugins/cache/dapper-claude-code-plugins/kaaos/2.0.0/skills/operational-rhythms/assets/launchd-configs/com.kaaos.weekly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Job identifier -->
    <key>Label</key>
    <string>com.kaaos.weekly</string>

    <!-- Human-readable description -->
    <key>ServiceDescription</key>
    <string>KAAOS Weekly Synthesis - Pattern extraction and playbook updates (~$4-5/run)</string>

    <!-- Command to execute -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/kaaos</string>
        <string>review</string>
        <string>weekly</string>
    </array>

    <!-- Schedule: Every Monday at 6:00 AM local time -->
    <!-- Weekday: 0=Sunday, 1=Monday, 2=Tuesday, etc. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/ben/.kaaos-knowledge</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>KAAOS_HOME</key>
        <string>/Users/ben/.kaaos-knowledge</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/Users/ben/.kaaos-knowledge/.kaaos/logs/weekly.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ben/.kaaos-knowledge/.kaaos/logs/weekly.error.log</string>

    <!-- Don't run at load, only on schedule -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Process configuration -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>2048</integer>
    </dict>

    <!-- Timeout: 15 minutes max for weekly synthesis -->
    <key>TimeOut</key>
    <integer>900</integer>

    <!-- Exit codes that indicate success -->
    <key>SuccessfulExitCodes</key>
    <array>
        <integer>0</integer>
    </array>

    <!-- Keep alive configuration -->
    <key>KeepAlive</key>
    <false/>

    <!-- Nice value (lower priority) -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- Prevent multiple instances -->
    <key>AbandonProcessGroup</key>
    <false/>
</dict>
</plist>
