# claude/.claude/plugins/cache/dapper-claude-code-plugins/kaaos/2.0.0/skills/operational-rhythms/assets/launchd-configs/com.kaaos.monthly.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <!-- Job identifier -->
    <key>Label</key>
    <string>com.kaaos.monthly</string>

    <!-- Human-readable description -->
    <key>ServiceDescription</key>
    <string>KAAOS Monthly Review - Strategic alignment with parallel agents (~$10-15/run)</string>

    <!-- Command to execute -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/kaaos</string>
        <string>review</string>
        <string>monthly</string>
    </array>

    <!-- Schedule: 1st of every month at 5:00 AM local time -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Working directory -->
    <key>WorkingDirectory</key>
    <string>/Users/ben/.kaaos-knowledge</string>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>KAAOS_HOME</key>
        <string>/Users/ben/.kaaos-knowledge</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <!-- Monthly uses parallel agents, may need more resources -->
        <key>KAAOS_MAX_PARALLEL_AGENTS</key>
        <string>3</string>
    </dict>

    <!-- Log output -->
    <key>StandardOutPath</key>
    <string>/Users/ben/.kaaos-knowledge/.kaaos/logs/monthly.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ben/.kaaos-knowledge/.kaaos/logs/monthly.error.log</string>

    <!-- Don't run at load, only on schedule -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Process configuration -->
    <key>ProcessType</key>
    <string>Background</string>

    <!-- Resource limits - higher for monthly review -->
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>4096</integer>
    </dict>

    <!-- Timeout: 30 minutes max for monthly review with parallel agents -->
    <key>TimeOut</key>
    <integer>1800</integer>

    <!-- Exit codes that indicate success -->
    <key>SuccessfulExitCodes</key>
    <array>
        <integer>0</integer>
    </array>

    <!-- Keep alive configuration -->
    <key>KeepAlive</key>
    <false/>

    <!-- Nice value (lower priority, runs in background) -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- Prevent multiple instances -->
    <key>AbandonProcessGroup</key>
    <false/>

    <!-- Low priority I/O -->
    <key>LowPriorityIO</key>
    <true/>
</dict>
</plist>
