# claude/.claude/plugins/cache/dapper-claude-code-plugins/kaaos/2.0.0/templates/config/launchd-example.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
=============================================================================
KAAOS Launchd Configuration Example
=============================================================================

This plist file demonstrates how to configure macOS launchd for KAAOS
automated operational rhythms.

USAGE:
  1. Copy this file to ~/Library/LaunchAgents/
  2. Rename to match the review type (e.g., com.kaaos.daily-review.plist)
  3. Edit the configuration below to match your setup
  4. Load the agent: launchctl load ~/Library/LaunchAgents/com.kaaos.daily-review.plist

SCHEDULE TYPES:
  - Daily Review:     Runs every morning at configured time
  - Weekly Synthesis: Runs on configured weekday
  - Monthly Review:   Runs on first day of each month
  - Quarterly Review: Runs on first day of Jan, Apr, Jul, Oct

IMPORTANT:
  - The script creates a marker file that KAAOS detects
  - Actual review runs when you next interact with Claude Code
  - This approach avoids running Claude in the background

For more examples, see: templates/scripts/README.md
=============================================================================
-->
<plist version="1.0">
<dict>
    <!-- =================================================================== -->
    <!-- Basic Configuration -->
    <!-- =================================================================== -->

    <!-- Unique identifier for this job -->
    <!-- Convention: com.kaaos.[review-type].plist -->
    <key>Label</key>
    <string>com.kaaos.daily-review</string>

    <!-- Human-readable description (optional but recommended) -->
    <key>Description</key>
    <string>KAAOS Daily Review Scheduler - Creates marker file for daily digest generation</string>

    <!-- =================================================================== -->
    <!-- Script Execution -->
    <!-- =================================================================== -->

    <!-- Path to the script to execute -->
    <key>ProgramArguments</key>
    <array>
        <!-- Shell to use -->
        <string>/bin/bash</string>
        <!-- Script path - UPDATE THIS to your actual path -->
        <string>/Users/YOUR_USERNAME/.kaaos/automation/daily-review.sh</string>
    </array>

    <!-- Working directory for the script -->
    <key>WorkingDirectory</key>
    <string>/Users/YOUR_USERNAME/kaaos-knowledge</string>

    <!-- =================================================================== -->
    <!-- Schedule Configuration -->
    <!-- =================================================================== -->

    <!--
    OPTION A: Calendar-based scheduling (recommended for reviews)

    Use StartCalendarInterval for specific times:
    - Minute: 0-59
    - Hour: 0-23 (24-hour format)
    - Day: 1-31 (day of month)
    - Weekday: 0-6 (0=Sunday, 1=Monday, etc.)
    - Month: 1-12

    Examples below show common configurations.
    -->

    <!-- DAILY REVIEW: Every day at 7:00 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!--
    WEEKLY SYNTHESIS: Every Monday at 6:00 AM
    Uncomment and use this instead for weekly:

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    -->

    <!--
    MONTHLY REVIEW: 1st of each month at 5:00 AM
    Uncomment and use this instead for monthly:

    <key>StartCalendarInterval</key>
    <dict>
        <key>Day</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>5</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    -->

    <!--
    QUARTERLY REVIEW: 1st of Jan, Apr, Jul, Oct at 3:00 AM
    For quarterly, create 4 separate plists or use multiple intervals:

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Month</key>
            <integer>1</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>3</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Month</key>
            <integer>4</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>3</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Month</key>
            <integer>7</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>3</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Month</key>
            <integer>10</integer>
            <key>Day</key>
            <integer>1</integer>
            <key>Hour</key>
            <integer>3</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>
    -->

    <!--
    OPTION B: Interval-based scheduling (for testing)

    Use StartInterval for recurring execution at fixed intervals (in seconds).
    Example: Run every 4 hours (14400 seconds)

    <key>StartInterval</key>
    <integer>14400</integer>
    -->

    <!-- =================================================================== -->
    <!-- Logging Configuration -->
    <!-- =================================================================== -->

    <!-- Standard output log -->
    <key>StandardOutPath</key>
    <string>/Users/YOUR_USERNAME/.kaaos/logs/daily-review.log</string>

    <!-- Standard error log -->
    <key>StandardErrorPath</key>
    <string>/Users/YOUR_USERNAME/.kaaos/logs/daily-review-error.log</string>

    <!-- =================================================================== -->
    <!-- Environment Variables -->
    <!-- =================================================================== -->

    <!-- Set environment variables for the script -->
    <key>EnvironmentVariables</key>
    <dict>
        <!-- KAAOS repository path -->
        <key>KAAOS_REPO</key>
        <string>/Users/YOUR_USERNAME/kaaos-knowledge</string>

        <!-- Review type identifier -->
        <key>KAAOS_REVIEW_TYPE</key>
        <string>daily</string>

        <!-- Optional: Notification settings -->
        <key>KAAOS_NOTIFY</key>
        <string>true</string>

        <!-- Optional: Debug mode -->
        <key>KAAOS_DEBUG</key>
        <string>false</string>
    </dict>

    <!-- =================================================================== -->
    <!-- Execution Options -->
    <!-- =================================================================== -->

    <!-- Run immediately if computer was asleep/off at scheduled time -->
    <key>RunAtLoad</key>
    <false/>

    <!-- Keep alive settings (not needed for KAAOS - script is quick) -->
    <key>KeepAlive</key>
    <false/>

    <!-- Launch only when user is logged in -->
    <key>LimitLoadToSessionType</key>
    <string>Aqua</string>

    <!-- Nice value for process priority (higher = lower priority) -->
    <key>Nice</key>
    <integer>10</integer>

    <!-- Process timeout in seconds (kill if takes too long) -->
    <key>TimeOut</key>
    <integer>60</integer>

    <!-- Throttle interval - minimum seconds between runs -->
    <key>ThrottleInterval</key>
    <integer>300</integer>

    <!-- =================================================================== -->
    <!-- Resource Limits -->
    <!-- =================================================================== -->

    <!-- Soft resource limits -->
    <key>SoftResourceLimits</key>
    <dict>
        <!-- Max memory in bytes (100 MB) -->
        <key>MemoryLock</key>
        <integer>104857600</integer>
    </dict>

</dict>
</plist>

<!--
=============================================================================
INSTALLATION INSTRUCTIONS
=============================================================================

1. CREATE YOUR REVIEW SCRIPTS
   Each review type needs a corresponding script in ~/.kaaos/automation/

   Example daily-review.sh:
   ```bash
   #!/bin/bash
   set -euo pipefail

   MARKER_DIR="$HOME/.kaaos/pending-reviews"
   mkdir -p "$MARKER_DIR"

   cat > "$MARKER_DIR/daily-$(date +%Y-%m-%d).marker" <<EOF
   type: daily
   created_at: $(date -Iseconds)
   status: pending
   EOF

   # Optional: Show notification
   osascript -e 'display notification "Daily digest ready" with title "KAAOS"'
   ```

2. COPY AND CONFIGURE THIS PLIST
   ```bash
   cp launchd-example.plist ~/Library/LaunchAgents/com.kaaos.daily-review.plist
   sed -i '' 's/YOUR_USERNAME/your_actual_username/g' ~/Library/LaunchAgents/com.kaaos.daily-review.plist
   ```

3. CREATE LOG DIRECTORY
   ```bash
   mkdir -p ~/.kaaos/logs
   ```

4. LOAD THE AGENT
   ```bash
   launchctl load ~/Library/LaunchAgents/com.kaaos.daily-review.plist
   ```

5. VERIFY IT'S RUNNING
   ```bash
   launchctl list | grep kaaos
   ```

6. CHECK NEXT SCHEDULED RUN
   ```bash
   launchctl print gui/$(id -u)/com.kaaos.daily-review
   ```

MANAGEMENT COMMANDS
===================

# List all KAAOS jobs
launchctl list | grep kaaos

# Unload (disable) a job
launchctl unload ~/Library/LaunchAgents/com.kaaos.daily-review.plist

# Reload after editing
launchctl unload ~/Library/LaunchAgents/com.kaaos.daily-review.plist
launchctl load ~/Library/LaunchAgents/com.kaaos.daily-review.plist

# Run immediately for testing
launchctl start com.kaaos.daily-review

# View logs
tail -f ~/.kaaos/logs/daily-review.log

# Check for errors
cat ~/.kaaos/logs/daily-review-error.log

TROUBLESHOOTING
===============

Problem: Job not running
- Check: launchctl list | grep kaaos (should show 0 exit status)
- Check: ls ~/Library/LaunchAgents/com.kaaos.*.plist
- Verify: plutil -lint ~/Library/LaunchAgents/com.kaaos.daily-review.plist

Problem: Script not found
- Ensure script path is absolute
- Check script has execute permission: chmod +x script.sh

Problem: Permission denied
- Ensure log directory exists and is writable
- Check script permissions

Problem: Running too often
- Increase ThrottleInterval
- Check StartCalendarInterval is correct

COMPLETE SETUP FOR ALL RHYTHMS
==============================

Create four plists and load them:

```bash
# Daily
cp launchd-example.plist ~/Library/LaunchAgents/com.kaaos.daily-review.plist
# Edit to use daily-review.sh and Hour=7

# Weekly
cp launchd-example.plist ~/Library/LaunchAgents/com.kaaos.weekly-synthesis.plist
# Edit to use weekly-synthesis.sh and Weekday=1, Hour=6

# Monthly
cp launchd-example.plist ~/Library/LaunchAgents/com.kaaos.monthly-review.plist
# Edit to use monthly-review.sh and Day=1, Hour=5

# Quarterly
cp launchd-example.plist ~/Library/LaunchAgents/com.kaaos.quarterly-review.plist
# Edit to use quarterly-review.sh and quarterly month array

# Load all
launchctl load ~/Library/LaunchAgents/com.kaaos.daily-review.plist
launchctl load ~/Library/LaunchAgents/com.kaaos.weekly-synthesis.plist
launchctl load ~/Library/LaunchAgents/com.kaaos.monthly-review.plist
launchctl load ~/Library/LaunchAgents/com.kaaos.quarterly-review.plist
```

=============================================================================
-->
